# housekeeping/etc/crontabs/root
# min   hour    day     month   weekday command
*       *       *       *       *       /bin/update-repose.sh
0       *       *       *       *       /bin/crashlogs-cleanup.sh
0       *       *       *       *       /bin/logrotate.sh
0       3       *       *       6       /bin/public-backup.sh
